# deploy/systemd/optus-hub-api.service
[Unit]
Description=Optus Hub API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=optushub
Group=optushub
WorkingDirectory=/opt/optus_hub/backend
EnvironmentFile=/opt/optus_hub/backend/.env
ExecStart=/opt/optus_hub/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
